# logviewer.socket
[Unit]
Description=Socket for LogViewer WEB Application
 
[Socket]
ListenStream=12943
 
[Install]
WantedBy=sockets.target